# data/kiran-system-daemon-biometrics.service.in
[Unit]
Description=Kiran System Daemon Biometrics Service

[Service]
Type=dbus
BusName=com.kylinsec.Kiran.SystemDaemon.Biometrics
ExecStart=@INSTALL_BINDIR@/kiran_biometrics_manager

[Install]
# We pull this in by graphical.target instead of waiting for the bus
# activation, to speed things up a little: gdm uses this anyway so it is nice
# if it is already around when gdm wants to use it and doesn't have to wait for
# it.
WantedBy=graphical.target
